# data/systemd/deepin-installer-cli.service
[Unit]
Description=Launch installer in normal mode

[Service]
Type=simple
ExecStart=/usr/bin/deepin-installer-cli-install
StandardInput=tty-force
TTYPath=/dev/tty2
TTYReset=yes
TTYVHangup=yes


